# examples/systemd/ocsp-responder.service
[Unit]
Description=OCSP Responder
Documentation=https://github.com/kuveris/ocsp
After=network.target

[Service]
Type=simple
User=ocsp
Group=ocsp
ExecStart=/usr/local/bin/ocsp-responder --config /etc/ocsp-responder/ocsp-responder.yaml
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Optional: systemd fails the unit (226/NAMESPACE) if a ReadWritePaths
# target is missing. Logs go to stdout and therefore the journal, so this
# is only needed if you redirect them to a file.
#ReadWritePaths=/var/log/ocsp-responder
ReadOnlyPaths=/etc/ocsp-responder /etc/ssl/certs

[Install]
WantedBy=multi-user.target
